# runner/launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
launchd.plist.example — пример launchd-агента для macOS.

Установка:
  1. Замени плейсхолдеры YOUR_USER и /full/path/to/ai-camp-info-helper
  2. Сохрани как ~/Library/LaunchAgents/com.aisurfers.info-helper.plist
  3. launchctl load ~/Library/LaunchAgents/com.aisurfers.info-helper.plist
  4. Проверь: launchctl list | grep info-helper

Расписание ниже — каждый будний день в 9:00.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aisurfers.info-helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/YOUR_USER/full/path/to/ai-camp-info-helper &amp;&amp; ./runner/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/full/path/to/ai-camp-info-helper/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/full/path/to/ai-camp-info-helper/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
